# supervisord.conf
[program:eurouteme]
command=gunicorn routeme:app -w 4 -b 0.0.0.0:80

[supervisord]
logfile=/home/ubuntu/routeme/supervisord.log
loglevel=debug
user=root
